# deployment/daphne.conf
[program:daphne]
command = /app/deployment/start_daphne.sh
directory=/app/core
user = root
autostart=true
autorestart=true
stdout_logfile = /logs/daphne.log
redirect_stderr = true