# deployment/daphne.conf
[program:daphne]
user=root
command = /app/deployment/start_daphne.sh
directory=/app/core
autostart=true
autorestart=true
stdout_logfile = /logs/daphne.log
redirect_stderr = true
